# packaging/acip-sidecar.service
[Unit]
Description=ACIP Sidecar
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=acip_user
Group=acip_user

# Recommended filesystem layout
WorkingDirectory=/opt/acip
ExecStart=/opt/acip/acip-sidecar \
  --config /etc/acip/config.toml \
  --secrets-file /etc/acip/secrets.env

Restart=on-failure
RestartSec=2

# Hardening (minimal; adjust as needed)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/acip /run/acip
RuntimeDirectory=acip
RuntimeDirectoryMode=0755
UMask=0007

# Environment
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
